# others/vulcan.conf
[program:vulcan]
command = /usr/bin/nodejs /home/ubuntu/www/vulcan/index.js 8000 
directory = /home/ubuntu/www/vulcan
user = root
